# backend/controlasist.service
[Unit]
Description=ControlAsist FastAPI Application
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/apizhe/backend
Environment=PATH=/var/www/apizhe/backend/venv/bin
ExecStart=/var/www/apizhe/backend/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 4
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# Configuración de seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/apizhe/backend/uploads /var/www/apizhe/backend/logs

# Configuración de logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=controlasist

[Install]
WantedBy=multi-user.target 